# gumstixFilesystem/etc/udev/rules.d/70-ide.rules
# IDE modules -- call script to load appropriate module
SUBSYSTEM=="ide", ACTION=="add", ATTRS{modalias}=="ide:*", RUN+="/sbin/modprobe $env{MODALIAS}"
